# 1_os_prep/playbooks/02-shared-fs/mount-qcow2-for-containerd.service
[Unit]
Description=Mount qcow2 image for containerd
Before=containerd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mount_qcow2_for_containerd.sh

[Install]
WantedBy=multi-user.target
